# supervisord.conf
[supervisord]
process_name=%(program_name)s_%(process_num)02d
nodaemon=true
user=root
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:php-fpm]
process_name=%(program_name)s_%(process_num)02d
command=/usr/sbin/php-fpm -F
stdout_events_enabled=true
user=root
stderr_events_enabled=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/project/artisan queue:work --queue=registration,funding,remittance,sendToAdmin,low --tries=250
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:httpd]
process_name=%(program_name)s_%(process_num)02d
command=/bin/bash -c "exec /usr/sbin/httpd -DFOREGROUND"
stdout_logfile=/dev/fd/1
user=root
stdout_logfile_maxbytes=0
redirect_stderr=true